# vsun.cron.d
0 2,14 * * * root /var/local/vsun/getvsun.sh

